# mosquitto/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs mosquitto
# ==============================================================================
declare -a options

options+=(-c /etc/mosquitto/mosquitto.conf)

if bashio::debug; then
  options+=(-v)
fi

# Send out discovery & service information
./discovery &

bashio::log.info "Starting mosquitto MQTT broker..."
exec mosquitto "${options[@]}"
